# etc/systemd/system/mysql.service.d/akura.conf
[Service]
LimitNOFILE=infinity
LimitMEMLOCK=infinity
OOMScoreAdjust=-1000
